# package/DEBIAN/prerm
#!/bin/bash
systemctl disable uav_configurator

